# crontab
0 1 * * * /usr/bin/php -q /home/www/sheet.php